# examples/launchd/com.chief-of-staff.today.plist
<!--
  Hourly "today" dispatcher (8:00 - 22:00)

  Setup:
    1. Copy this file to ~/Library/LaunchAgents/
    2. Replace ALL occurrences of YOUR_HOME with your actual home path (e.g., /Users/yourname)
    3. Replace ALL occurrences of YOUR_WORKSPACE with your workspace path (e.g., /Users/yourname/workspace)
    4. Load: launchctl load ~/Library/LaunchAgents/com.chief-of-staff.today.plist
    5. Unload: launchctl unload ~/Library/LaunchAgents/com.chief-of-staff.today.plist
    6. Check status: launchctl list | grep chief-of-staff

  Note: launchd does NOT expand $HOME or ~ in plist values. You must use absolute paths.
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chief-of-staff.today</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>YOUR_WORKSPACE/scripts/autonomous/dispatcher.sh</string>
        <string>today</string>
    </array>

    <key>WorkingDirectory</key>
    <string>YOUR_WORKSPACE</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>17</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>19</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/chief-of-staff-today.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/chief-of-staff-today.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>YOUR_HOME</string>
    </dict>
</dict>
</plist>
